# ansible/roles/webserver/files/rc.local
#!/bin/sh -e


/usr/local/bin/uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
exit 0